# res/openrc-service
#!/sbin/openrc-run
supervisor=supervise-daemon
description="emptty Display Manager"
command=/usr/bin/emptty
command_args="-d"
pidfile="/run/${RC_SVCNAME}.pid"
respawn_period=${respawn_period-60}
